# software/configs/hotspot-udev.rules
SUBSYSTEM=="net", ENV{DEVTYPE}=="wlan", ACTION=="add", RUN+="/opt/apertus/axiom-firmware/software/configs/create_hotspot.sh"
